# caddy-cert-watcher.path
[Unit]
Description=Watch for Caddy certificate renewals
Documentation=https://caddyserver.com/docs/

[Path]
PathModified=/var/lib/caddy/.local/share/caddy/certificates/acme-v02.api.letsencrypt.org-directory
Unit=caddy-cert-watcher.service

[Install]
WantedBy=multi-user.target
